# base/univention-server/conffiles/etc/cron.d/univention-managed-client
@%@UCRWARNING=# @%@

PATH=/usr/sbin:/usr/bin:/sbin:/bin
@%@ldap/connection/check@%@ root /usr/share/univention-managed-client/check_connection "runParts" >/var/log/univention/check_connection.log 2>&1
@reboot root /usr/share/univention-managed-client/check_connection "runParts" >/var/log/univention/check_connection.log 2>&1
