# base/univention-server/conffiles/etc/cron.d/univention-server-backup
PATH=/usr/sbin:/usr/bin:/sbin:/bin
 0 1 * * * root /usr/sbin/jitter 600 /usr/lib/univention-server/server_password_change
