# base/univention-server/conffiles/etc/cron.d/univention-mobile-client
PATH=/usr/sbin:/usr/bin:/sbin:/bin
@%@ldap/connection/check@%@ root /usr/share/univention-mobile-client/check_connection
@reboot root /usr/share/univention-mobile-client/check_connection
