# mail/univention-postgrey/debian/postinst
#!/bin/sh

univention-config-registry set "postfix/greylisting?enabled"
univention-config-registry set "postfix/greylisting/delay?300"
univention-config-registry set "postfix/greylisting/max-age?35"
univention-config-registry set "postfix/greylisting/retry-window?48"
univention-config-registry set "postfix/greylisting/lookup?host"
univention-config-registry set "postfix/greylisting/privacy?true"
univention-config-registry set "postfix/greylisting/client/whitelist/auto?5"
univention-config-registry set "postfix/greylisting/client/whitelist?/etc/postgrey/whitelist_clients /etc/postgrey/whitelist_clients.local"
univention-config-registry set "postfix/greylisting/recipient/whitelist?/etc/postgrey/whitelist_recipients /etc/postgrey/whitelist_recipients.local"

univention-config-registry commit "/etc/init.d/postgrey"

if [ ! -f "/etc/postgrey/whitelist_recipients.local" ]
then
	install -m 0644 -o 0 -g 0 /dev/null "/etc/postgrey/whitelist_recipients.local"
fi
if [ ! -f "/etc/postgrey/whitelist_clients.local" ]
then
	install -m 0644 -o 0 -g 0 /dev/null "/etc/postgrey/whitelist_clients.local"
fi

invoke-rc.d postgrey restart
invoke-rc.d postfix restart

exit 0
