# services/univention-bind/conffiles/etc/runit/univention-bind/run
#!/bin/sh

# Warning: This file is auto-generated and might be overwritten by 
#          univention-config-registry.
#          Please edit the following file instead:
# Warnung: Diese Datei wurde automatisch generiert und kann durch 
#          univention-config-registry überschrieben werden. 
#          Bitte bearbeiten Sie an Stelle dessen die folgende Datei:
# 
# 	/etc/univention/templates/files/etc/runit/univention-bind/run
# 


PATH=/sbin:/bin:/usr/sbin:/usr/bin

OPTS="-c /etc/bind/named.conf -p @%@dns/master/port@%@ -u bind -f"

# check ucr autostart setting
if [ -f "/usr/share/univention-config-registry/init-autostart.lib" ]; then
    source "/usr/share/univention-config-registry/init-autostart.lib"
    check_autostart univention-bind bind/autostart
fi

exec /usr/sbin/named $OPTS
