# services/univention-spamassassin/debian/univention-spamassassin.conffiles
/etc/univention/templates/files/usr/share/spamassassin/10_misc.cf
/etc/univention/templates/files/etc/default/spamassassin

/etc/univention/templates/files/usr/sbin/univention-sa-learn
/etc/univention/templates/files/etc/cron.daily/univention-spamassassin
