# services/univention-bind/conffiles/etc/init.d/bind9
#!/bin/sh
@%@UCRWARNING=# @%@
### BEGIN INIT INFO
# Provides:          bind9
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
@!@
SERVICES = {
	'samba4': 'samba-ad-dc',
	'ldap': 'slapd',
}
backend = configRegistry.get('dns/backend', 'ldap')
service = SERVICES.get(backend, '')
print("# Required-Start:    %s" % (service,))
print("# Required-Stop:     %s" % (service,))
@!@# Short-Description: bind9 Domain Name Server (DNS)
### END INIT INFO
#
# Copyright 2001-2019 Univention GmbH
#
# https://www.univention.de/
#
# All rights reserved.
#
# The source code of this program is made available
# under the terms of the GNU Affero General Public License version 3
# (GNU AGPL V3) as published by the Free Software Foundation.
#
# Binary versions of this program provided by Univention to you as
# well as other copyrighted, protected or trademarked materials like
# Logos, graphics, fonts, specific documentations and configurations,
# cryptographic keys etc. are subject to a license agreement between
# you and Univention and not subject to the GNU AGPL V3.
#
# In the case you use this program under the terms of the GNU AGPL V3,
# the program is provided in the hope that it will be useful,
# but WITHOUT ANY WARRANTY; without even the implied warranty of
# MERCHANTABILITY or FITNESS FOR A PARTICULAR PURPOSE. See the
# GNU Affero General Public License for more details.
#
# You should have received a copy of the GNU Affero General Public
# License with the Debian GNU/Linux or Univention distribution in file
# /usr/share/common-licenses/AGPL-3; if not, see
# <https://www.gnu.org/licenses/>.
#

# sourcing /lib/lsb/init-functions calls bind9.service if it exists
. /lib/lsb/init-functions

# During systemsetup inside chroot we need a simple init script

eval "$(univention-config-registry shell dns/backend)"

NAME='bind9'
DESC=$(get_lsb_header_val "$0" Short-Description)

case "$1" in
	start)
		case "$dns_backend" in
			samba4)
				start-stop-daemon --background --start --exec /usr/lib/univention-bind/samba4 -- start
				/usr/lib/univention-bind/samba4 wait-for-startup
				;;
			ldap)
				start-stop-daemon --background --start --exec /usr/lib/univention-bind/proxy -- start
				start-stop-daemon --background --start --exec /usr/lib/univention-bind/ldap -- start
				/usr/lib/univention-bind/proxy wait-for-startup
				;;
		esac
		;;
	stop)
		case "$dns_backend" in
			samba4)
				/usr/lib/univention-bind/samba4 stop
				;;
			ldap)
				/usr/lib/univention-bind/proxy stop
				/usr/lib/univention-bind/ldap stop
				;;
		esac
		;;
	restart)
		case "$dns_backend" in
			samba4)
				/usr/lib/univention-bind/samba4 stop
				start-stop-daemon --background --start --exec /usr/lib/univention-bind/samba4 -- start
				/usr/lib/univention-bind/samba4 wait-for-startup
				;;
			ldap)
				/usr/lib/univention-bind/proxy stop
				/usr/lib/univention-bind/ldap stop
				start-stop-daemon --background --start --exec /usr/lib/univention-bind/proxy -- start
				start-stop-daemon --background --start --exec /usr/lib/univention-bind/ldap -- start
				/usr/lib/univention-bind/proxy wait-for-startup
				;;
		esac
		;;
	reload)
		case "$dns_backend" in
			samba4)
				/usr/lib/univention-bind/samba4 reload
				;;
			ldap)
				/usr/lib/univention-bind/proxy reload
				/usr/lib/univention-bind/ldap reload
				;;
		esac
		;;

esac
exit 0
