# services/univention-cloud-init/etc/init.d/univention-cloud-init-prepare-boot
#!/bin/bash
### BEGIN INIT INFO
# Provides:          univention-cloud-init-prepare-boot
# Required-Start:    $network $remote_fs
# Required-Stop:     $network $remote_fs
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Univention Cloud initialization boot preparation
### END INIT INFO
#
# Copyright 2014-2019 Univention GmbH
#
# https://www.univention.de/
#
# All rights reserved.
#
# The source code of this program is made available
# under the terms of the GNU Affero General Public License version 3
# (GNU AGPL V3) as published by the Free Software Foundation.
#
# Binary versions of this program provided by Univention to you as
# well as other copyrighted, protected or trademarked materials like
# Logos, graphics, fonts, specific documentations and configurations,
# cryptographic keys etc. are subject to a license agreement between
# you and Univention and not subject to the GNU AGPL V3.
#
# In the case you use this program under the terms of the GNU AGPL V3,
# the program is provided in the hope that it will be useful,
# but WITHOUT ANY WARRANTY; without even the implied warranty of
# MERCHANTABILITY or FITNESS FOR A PARTICULAR PURPOSE. See the
# GNU Affero General Public License for more details.
#
# You should have received a copy of the GNU Affero General Public
# License with the Debian GNU/Linux or Univention distribution in file
# /usr/share/common-licenses/AGPL-3; if not, see
# <https://www.gnu.org/licenses/>.

. /lib/lsb/init-functions

. /usr/share/univention-lib/ucr.sh
NAME='univention-cloud-init-prepare-boot'
DESC=$(get_lsb_header_val "$0" Short-Description)

case "$1" in
	start)
		log_action_msg "Running $DESC"
		if is_ucr_true server/cloud-init; then
			ucr commit /var/www/univention/meta.json >/dev/null 2>&1
			ucr set apache2/startsite="univention/cloud-init-setup.html"
		fi
		log_action_end_msg 0
		;;
	stop|restart|force-reload)
		;;
	status)
		if is_ucr_true server/cloud-init
		then
			log_success_msg "$DESC is enabled"
		else
			log_failure_msg "$DESC is disabled"
		fi
		;;
	*)
		echo "Usage: $0 {start|stop|restart|force-reload|status}."
		exit 1
		;;
esac
