# management/univention-management-console/systemd/system/univention-management-console-server-multiprocessing.target
[Unit]
Description=Start UMC-Server in multiprocessing mode
PartOf=univention-management-console-server.service

[Install]
WantedBy=multi-user.target
