# management/univention-management-console/systemd/system/univention-management-console-web-server-multiprocessing.target
[Unit]
Description=Start UMC-Web-Server in multiprocessing mode
PartOf=univention-management-console-web-server.service

[Install]
WantedBy=multi-user.target
