# management/univention-management-console/systemd/system/univention-management-console-web-server.service.d/dependencies.conf
[Unit]
Wants=univention-management-console-web-server-multiprocessing.target
