# services/univention-nfs/conffiles/etc/init.d/nfs-kernel-server
#!/bin/sh
#
# Univention NFS
#  init script
#
# Copyright 2004-2012 Univention GmbH
#
# http://www.univention.de/
#
# All rights reserved.
#
# The source code of this program is made available
# under the terms of the GNU Affero General Public License version 3
# (GNU AGPL V3) as published by the Free Software Foundation.
#
# Binary versions of this program provided by Univention to you as
# well as other copyrighted, protected or trademarked materials like
# Logos, graphics, fonts, specific documentations and configurations,
# cryptographic keys etc. are subject to a license agreement between
# you and Univention and not subject to the GNU AGPL V3.
#
# In the case you use this program under the terms of the GNU AGPL V3,
# the program is provided in the hope that it will be useful,
# but WITHOUT ANY WARRANTY; without even the implied warranty of
# MERCHANTABILITY or FITNESS FOR A PARTICULAR PURPOSE. See the
# GNU Affero General Public License for more details.
#
# You should have received a copy of the GNU Affero General Public
# License with the Debian GNU/Linux or Univention distribution in file
# /usr/share/common-licenses/AGPL-3; if not, see
# <http://www.gnu.org/licenses/>.
#
# nfs-kernel-server
#		This shell script takes care of starting and stopping
#               the kernel-mode NFS server.
#
# chkconfig: 345 60 20
# description: NFS is a popular protocol for file sharing across TCP/IP \
#              networks. This service provides NFS server functionality, \
#              which is configured via the /etc/exports file.
#
@%@UCRWARNING=# @%@

PREFIX=/usr
[ -x $PREFIX/sbin/rpc.nfsd ] || exit 0
[ -x $PREFIX/sbin/rpc.mountd ] || exit 0
[ -x $PREFIX/sbin/exportfs ] || exit 0

RPCNFSDCOUNT=8	# Number of servers to be started up by default
RPCMOUNTDOPTS=

. /lib/lsb/init-functions

# If we should be started by heartbeat exit here
eval "$(/usr/sbin/univention-config-registry shell)"

# See how we were called.
case "$1" in
  start)
	# check ucr autostart setting
	if [ -f "/usr/share/univention-config-registry/init-autostart.lib" ]; then
		. /usr/share/univention-config-registry/init-autostart.lib
		check_autostart nfs-kernel-server nfs/autostart
	fi
	cd /	# daemons should have root dir as cwd
	if egrep -q '^["]{0,1}/' /etc/exports
	then
		log_action_msg "Exporting directories for NFS kernel daemon"
		$PREFIX/sbin/exportfs -r
		echo "done."

		log_action_msg "Starting NFS kernel daemon: nfsd"
		start-stop-daemon --start --quiet \
		    --exec $PREFIX/sbin/rpc.nfsd -- $RPCNFSDCOUNT
		log_action_msg "Starting NFS kernel daemon: mountd"
		$PREFIX/bin/rpcinfo -u localhost nfs 3 >/dev/null 2>&1 ||
		    RPCMOUNTDOPTS="$RPCMOUNTDOPTS --no-nfs-version 3"
		start-stop-daemon --start --quiet \
		    --exec $PREFIX/sbin/rpc.mountd -- $RPCMOUNTDOPTS
		log_action_end_msg 0
	else
		log_action_msg "Not starting NFS kernel daemon: No exports."
		log_action_end_msg 1
	fi
	;;

  stop)
	log_action_msg "Stopping NFS kernel daemon: mountd"
	start-stop-daemon --stop --oknodo --quiet \
	    --name rpc.mountd --user 0
	log_action_msg "Stopping NFS kernel daemon: nfsd"
	start-stop-daemon --stop --oknodo --quiet \
	    --name nfsd --user 0 --signal 2
	log_action_end_msg 0

	log_action_msg "Unexporting directories for NFS kernel daemon"
	$PREFIX/sbin/exportfs -au
	log_action_end_msg 0
	;;

  reload | force-reload)
	log_action_msg "Re-exporting directories for NFS kernel daemon"
	$PREFIX/sbin/exportfs -r
	log_action_end_msg 0
	;;

  restart)
	$0 stop
	sleep 1
	$0 start
	;;

  *)
	echo "Usage: nfs-kernel-server {start|stop|reload|force-reload|restart}"
	exit 1
	;;
esac

exit 0
