# base/univention-installer/etc/inittab
# /etc/inittab

::sysinit:/sbin/univention-installer-startup

#installer
tty1::respawn:/sbin/univention-installer

tty2::askfirst:/bin/sh
tty3::askfirst:/bin/sh

#tty4::respawn:/bin/tail -f /var/log/messages
tty4::respawn:/bin/tail -f /var/log/syslogd


::ctrlaltdel:/bin/reboot
