# services/univention-mariadb/conffiles/lib/systemd/system/mariadb.service.d/univention.conf
@%@UCRWARNING=# @%@

[Service]
@!@
print ('TimeoutStartSec=%s' % configRegistry.get('mariadb/startup/timeout', '90'))
@!@
