# services/univention-samba/debian/postinst
#! /bin/sh
#
# Univention Samba
#  postinst script for the debian package
#
# Copyright 2001-2012 Univention GmbH
#
# http://www.univention.de/
#
# All rights reserved.
#
# The source code of this program is made available
# under the terms of the GNU Affero General Public License version 3
# (GNU AGPL V3) as published by the Free Software Foundation.
#
# Binary versions of this program provided by Univention to you as
# well as other copyrighted, protected or trademarked materials like
# Logos, graphics, fonts, specific documentations and configurations,
# cryptographic keys etc. are subject to a license agreement between
# you and Univention and not subject to the GNU AGPL V3.
#
# In the case you use this program under the terms of the GNU AGPL V3,
# the program is provided in the hope that it will be useful,
# but WITHOUT ANY WARRANTY; without even the implied warranty of
# MERCHANTABILITY or FITNESS FOR A PARTICULAR PURPOSE. See the
# GNU Affero General Public License for more details.
#
# You should have received a copy of the GNU Affero General Public
# License with the Debian GNU/Linux or Univention distribution in file
# /usr/share/common-licenses/AGPL-3; if not, see
# <http://www.gnu.org/licenses/>.

. /usr/share/univention-lib/all.sh

eval "$(univention-config-registry shell)"

create_logfile /var/log/univention/samba-sync.log "root:adm" 640
create_logfile /var/log/univention/pdc-takeover.log "root:adm" 640

if [ -z "$2" ]; then
	mkdir -p /home/groups

	mkdir -p /var/lib/samba/netlogon/scripts
	if [ ! -e /var/lib/samba/netlogon/scripts/startup.cmd ]; then
		cp /usr/share/doc/univention-samba/startup.cmd /var/lib/samba/netlogon/scripts/startup.cmd;
	fi;

	mkdir -p /etc/univention/skel/windows-profiles
	for dname in Win95 WinNT Win2k WinXP Win2K3 Vista Vista.V2 ; do
		mkdir -p /etc/univention/skel/windows-profiles/$dname
	done
fi

univention-config-registry set samba/share/home?yes \
							   samba/share/groups?no \
							   samba/adminusers?"administrator join-backup"

if [ "$server_role" != "memberserver" ]
then
	if [ -n "$samba_netbios_name" ]
	then
		tmphostname=$samba_netbios_name
	elif [ -n "$samba_ha_master" ]
	then
		tmphostname="$samba_ha_master"
	else
		tmphostname="$hostname"
	fi
	univention-config-registry set samba/profileserver?$tmphostname \
								   samba/profilepath?'%U\windows-profiles\%a' \
								   samba/homedirserver?$tmphostname \
								   samba/homedirpath?%U \
								   samba/homedirletter?I
fi

#DEBHELPER#

if [ -z "$windows_domain" ]; then
	domain=`echo $domainname | cut -d"." -f 1 | tr "[:lower:]" "[:upper:]"`
	univention-config-registry set windows/domain=$domain
fi

univention-config-registry set samba/debug/level?0 \
							   samba/os/level?65

if [ "$1" = "configure" ] && [ -z "$2" ]; then
	# only set this for new installations
	univention-config-registry set samba/autostart="no"
fi

call_joinscript 26univention-samba.inst

univention-config-registry set samba/script/adduser?true \
				samba/script/deleteuser?true \
				samba/script/addgroup?true \
				samba/script/deletegroup?true \
				samba/script/addusertogroup?true \
				samba/script/deleteuserfromgroup?true \
				samba/script/addmachine?true \
				samba/script/setprimarygroup?true \
				samba/script/postusermodify?false \
				samba/winbind/nested/groups?no \
				samba/encrypt_passwords?yes \
				samba/use_spnego?yes \
				samba/client_use_spnego?yes \
				samba/oplocks?yes \
				samba/kernel_oplocks?yes \
				samba/large_readwrite?yes \
				samba/deadtime?15 \
				samba/read_raw?yes \
				samba/write_raw?yes \
				samba/max_xmit?65535 \
				samba/max_open_files?32808 \
				samba/getwd_cache?yes \
				samba/store_dos_attributes?yes \
				samba/preserve_case?yes \
				samba/short_preserve_case?yes \
				samba/time_server?yes \
				samba/guest_account?nobody \
				samba/map_to_guest?"Bad User" \
				samba/netlogon/sync?sync \
				samba/domain/logons?auto \
				samba/password/checkscript?"/usr/share/univention-samba/password_check %u"

if [ "$1" = "configure" ]; then

	if [ -z "$2" ]; then
		# only set this for new installations
		if [ "$server_role" = "domaincontroller_master" ]; then
			univention-config-registry set windows/wins-support=yes
		else
			univention-config-registry set windows/wins-support?no
			if [ -n "$ldap_master" ]; then
				univention-config-registry set windows/wins-server?$ldap_master
			fi
		fi
	elif [ -n "$2" ]; then
		/etc/init.d/samba crestart
	fi


	if test -f /etc/init.d/univention-directory-listener; then
		/etc/init.d/univention-directory-listener crestart
	fi
fi

# configure firewall
ucr set security/packetfilter/package/univention-samba/tcp/137:139/all="ACCEPT" \
		security/packetfilter/package/univention-samba/tcp/137:139/all/en="netbios (Samba)" \
		security/packetfilter/package/univention-samba/udp/137:139/all="ACCEPT" \
		security/packetfilter/package/univention-samba/udp/137:139/all/en="netbios (Samba)" \
		security/packetfilter/package/univention-samba/udp/137/all="ACCEPT" \
		security/packetfilter/package/univention-samba/tcp/445/all="ACCEPT" \
		security/packetfilter/package/univention-samba/tcp/445/all/en="microsoft-ds (Samba)" \
		security/packetfilter/package/univention-samba/udp/445/all="ACCEPT" \
		security/packetfilter/package/univention-samba/udp/445/all/en="microsoft-ds (Samba)"
[ -x "/etc/init.d/univention-firewall" ] && invoke-rc.d univention-firewall restart

exit 0

