# mail/univention-mail-postfix/lib/systemd/system/univention-postfix-sender-check-smtp.service
[Unit]
Description=Postfix milter to prevent forged From (SMTP)
After=slapd.service

[Service]
Type=simple
ExecStart=/usr/share/univention-mail-postfix/sender_check_milter_smtp
ExecReload=/bin/kill -s HUP $MAINPID
Restart=always
RestartSec=3
User=listfilter
Group=nogroup

[Install]
WantedBy=multi-user.target
