# base/univention-home-mounter/conffiles/etc/cron.d/univention-home-mounter
@%@UCRWARNING=# @%@
PATH=/usr/sbin:/usr/bin:/sbin:/bin

# /etc/cron.d/univention-home-mounter: crontab fragment to unmount home directories

# Run queue every 10 minutes
*/10 * * * *     root   if [ -x /usr/sbin/univention-umount-homedirs ]; then /usr/sbin/univention-umount-homedirs; fi
