# management/univention-self-service/conffiles/etc/supervisor/conf.d/self-service.conf
@%@UCRWARNING=# @%@

[program:self-service]
command=/usr/bin/python /usr/share/univention-self-service/self-service-web.py
directory=/usr/share/univention-self-service/
user=self-service

