# base/univention-maintenance-mode/univention-maintenance-mode.service
[Unit]
Description=A UCS maintenance mode site
After=network.target

[Service]
Type=simple
ExecStart=/usr/share/univention-maintenance-mode-update-progress
