# virtualization/univention-virtual-machine-manager-node/conffiles/etc/runit/univention-libvirt/run
#!/bin/bash
set -e

PATH=/sbin:/bin:/usr/sbin:/usr/bin

source /etc/default/libvirt-bin

exec /usr/sbin/libvirtd ${libvirtd_opts/-d/}
