# base/univention-bootsplash/systemd-units/univention-welcome-screen-daemon.service
[Unit]
Description=Starting univention-welcome screen daemon

[Service]
ExecStart=/usr/bin/univention-welcome-screen
