# management/univention-management-console/conffiles/etc/runit/univention-management-console-proxy/run
#!/bin/sh

# Warning: This file is auto-generated and might be overwritten by
#          univention-baseconfig.
#          Please edit the following file instead:
# Warnung: Diese Datei wurde automatisch generiert und kann durch
#          univention-baseconfig überschrieben werden.
#          Bitte bearbeiten Sie an Stelle dessen die folgende Datei:
#
# 	/etc/univention/templates/files/etc/runit/univention-management-console-proxy/run
#


PATH=/sbin:/bin:/usr/sbin:/usr/bin

exec /usr/sbin/univention-management-console-proxy -n
