# management/univention-self-service/debian/univention-self-service-passwordreset-umc.service
[Unit]
Description=memcached daemon (univention-self-service instance)
After=network.target

[Service]
ExecStart=/usr/share/memcached/scripts/systemd-memcached-wrapper /etc/memcached_univention-self-service.conf

[Install]
WantedBy=multi-user.target
