# services/univention-bind/lib/systemd/system/bind9.service
[Unit]
Documentation=man:named(8)
Wants=nss-lookup.target
Before=nss-lookup.target

[Service]
Restart=on-failure

[Install]
WantedBy=multi-user.target
