# base/univention-updater/debian/postinst
#! /bin/sh
#
# Univention Updater
#  postinst script
#
# Copyright (C) 2004, 2005, 2006, 2007 Univention GmbH
#
# http://www.univention.de/
#
# All rights reserved.
#
# This program is free software; you can redistribute it and/or modify
# it under the terms of the GNU General Public License version 2 as
# published by the Free Software Foundation.
#
# Binary versions of this file provided by Univention to you as
# well as other copyrighted, protected or trademarked materials like
# Logos, graphics, fonts, specific documentations and configurations,
# cryptographic keys etc. are subject to a license agreement between
# you and Univention.
#
# This program is distributed in the hope that it will be useful,
# but WITHOUT ANY WARRANTY; without even the implied warranty of
# MERCHANTABILITY or FITNESS FOR A PARTICULAR PURPOSE.  See the
# GNU General Public License for more details.
#
# You should have received a copy of the GNU General Public License
# along with this program; if not, write to the Free Software
# Foundation, Inc., 51 Franklin St, Fifth Floor, Boston, MA  02110-1301  USA

createLogfile () {
	if [ ! -e $1 ] ; then
		touch $1
		chown $2 $1
		chmod $3 $1
	fi
}

createLogfile /var/log/univention/repository.log "root:adm" 640
createLogfile /var/log/univention/security-updates.log "root:adm" 640
createLogfile /var/log/univention/updater.log "root:adm" 640
createLogfile /var/log/univention/actualise.log "root:adm" 640

if [ "$1" = "configure" ]; then
	eval `univention-baseconfig shell version/version`
	eval `univention-baseconfig shell version/patchlevel`
	if [ -z "$version_version" -o -z "$version_patchlevel" ]; then
		echo "Error: baseconf-variables version/version and/or version/patchlevel are not set."
		exit 1
	fi

	eval `univention-baseconfig shell server/role`
	if [ "$server_role" = "domaincontroller_master" ]; then
		univention-baseconfig set update/server?"http://download.univention.de"
	else
		univention-baseconfig set update/server?"http://univention-repository/univention-cdrom"
	fi

	univention-baseconfig set update/security/server?"http://download.univention.de"

	eval `univention-baseconfig shell update/security/server`
	eval `univention-baseconfig shell update/server`

	if [ "$update_server" == "http://bitz150.bitz.briteline.de" ]; then
		univention-baseconfig set update/server="http://download.univention.de"
	fi

	if [ "$update_security_server" == "http://bitz150.bitz.briteline.de" ]; then
		univention-baseconfig set update/security/server="http://download.univention.de"
	fi

	univention-baseconfig set update/security/prefix="updates/security-desktop"
	univention-baseconfig set update/directory?"ucs-updates"

	mkdir -p /var/cache/univention-update

	if [ -n "$2" ] && dpkg --compare-versions "$2" lt 1.12; then
		# copy sources.list
		if [ -e "/etc/apt/sources.list" ]; then
			cp /etc/apt/sources.list /etc/apt/sources.list.ucs1_2_3
		fi
		eval `univention-baseconfig shell`
		if [ -d "/var/lib/univention-server-cdrom" ]; then
			univention-baseconfig set local/repository=true
			touch /var/lib/univention-server-cdrom/sources.list
			if [ -d "/var/lib/univention-server-cdrom/packages" ]; then
				echo "deb http://$hostname.$domainname/univention-cdrom/packages/ ./">> /var/lib/univention-server-cdrom/sources.list
			fi
			mkdir -p /var/lib/univention-server-cdrom/ucs-updates
			touch /var/lib/univention-server-cdrom/ucs-updates/index.list

			if [ -e /var/lib/univention-server-cdrom/.univention_install ]; then
				vv=`cat /var/lib/univention-server-cdrom/.univention_install | grep "VERSION=" | sed -e 's|VERSION=||g'`
				vp=`cat /var/lib/univention-server-cdrom/.univention_install | grep "PATCHLEVEL=" | sed -e 's|PATCHLEVEL=||g'`
				echo "$vv-$vp" >>/var/lib/univention-server-cdrom/ucs-updates/index.list
			fi

			for dir in /var/lib/univention-server-cdrom/ucs-updates/*; do
				if [ -d "$dir" ]; then
					echo "deb http://$hostname.$domainname/univention-cdrom/ucs-updates/`basename $dir`/ ./">> /var/lib/univention-server-cdrom/sources.list
				fi
			done

			first=1

			for file in /var/lib/univention-server-cdrom/ucs-updates/*.txt; do
				if [ $first = 1 ]; then
					echo "">/var/lib/univention-server-cdrom/ucs-updates/index.list
					first=0
				fi
				echo `basename $file` | sed -e 's|.txt$||;s|^ucs_||' >>/var/lib/univention-server-cdrom/ucs-updates/index.list
			done
		fi
	fi
	# update cronjob-entries
	if [ -x /usr/lib/univention-directory-policy/univention-policy-maintenance ]; then
		/usr/lib/univention-directory-policy/univention-policy-maintenance
	fi
	# coloured warning in univention-maintenance
	univention-baseconfig set update/warning?"no"
	univention-baseconfig set update/warning/coloured?"no"
	univention-baseconfig set update/warning/lang?"default"
	univention-baseconfig set update/warning/tty?"/dev/tty1 /dev/tty2 /dev/tty3 /dev/tty4 /dev/tty5 /dev/tty6 /dev/tty7"

fi

if [ -n "$2" ] && dpkg --compare-versions "$2" lt 1.36-7; then
	#create index.list
	if [ -d "/var/lib/univention-server-cdrom/ucs-updates" ]; then
		rm /var/lib/univention-server-cdrom/ucs-updates/index.list
		for txt in /var/lib/univention-server-cdrom/ucs-updates/*.txt; do
			b_txt=`basename $txt`
			version=`echo "$b_txt" | sed -e 's|ucs_||;s|.txt||'`
			echo "$version" >>/var/lib/univention-server-cdrom/ucs-updates/index.list
		done
	fi
fi

if [ -n "$2" ] && dpkg --compare-versions "$2" lt 1.36-42; then
	update-rc.d -f univention-maintenance remove || true
fi

#DEBHELPER#

exit 0


