# base/univention-updater/debian/postinst
#! /bin/sh
#
# Univention Updater
#  postinst script
#
# Copyright (C) 2004, 2005, 2006, 2007, 2008 Univention GmbH
#
# http://www.univention.de/
#
# All rights reserved.
#
# This program is free software; you can redistribute it and/or modify
# it under the terms of the GNU General Public License version 2 as
# published by the Free Software Foundation.
#
# Binary versions of this file provided by Univention to you as
# well as other copyrighted, protected or trademarked materials like
# Logos, graphics, fonts, specific documentations and configurations,
# cryptographic keys etc. are subject to a license agreement between
# you and Univention.
#
# This program is distributed in the hope that it will be useful,
# but WITHOUT ANY WARRANTY; without even the implied warranty of
# MERCHANTABILITY or FITNESS FOR A PARTICULAR PURPOSE.  See the
# GNU General Public License for more details.
#
# You should have received a copy of the GNU General Public License
# along with this program; if not, write to the Free Software
# Foundation, Inc., 51 Franklin St, Fifth Floor, Boston, MA  02110-1301  USA

createLogfile () {
	if [ ! -e $1 ] ; then
		touch $1
		chown $2 $1
		chmod $3 $1
	fi
}

createLogfile /var/log/univention/repository.log "root:adm" 640
createLogfile /var/log/univention/security-updates.log "root:adm" 640
createLogfile /var/log/univention/updater.log "root:adm" 640
createLogfile /var/log/univention/actualise.log "root:adm" 640

if [ "$1" = "configure" ]; then
	# update cronjob-entries
	if [ -x /usr/lib/univention-directory-policy/univention-policy-maintenance ]; then
		/usr/lib/univention-directory-policy/univention-policy-maintenance
	fi
	# coloured warning in univention-maintenance
	univention-baseconfig set update/warning?"no"
	univention-baseconfig set update/warning/coloured?"no"
	univention-baseconfig set update/warning/lang?"default"
	univention-baseconfig set update/warning/tty?"/dev/tty1 /dev/tty2 /dev/tty3 /dev/tty4 /dev/tty5 /dev/tty6 /dev/tty7"

fi

if [ -n "$2" ] && dpkg --compare-versions "$2" lt 1.36-7; then
	#create index.list
	if [ -d "/var/lib/univention-server-cdrom/ucs-updates" ]; then
		rm /var/lib/univention-server-cdrom/ucs-updates/index.list
		for txt in /var/lib/univention-server-cdrom/ucs-updates/*.txt; do
			b_txt=`basename $txt`
			version=`echo "$b_txt" | sed -e 's|ucs_||;s|.txt||'`
			echo "$version" >>/var/lib/univention-server-cdrom/ucs-updates/index.list
		done
	fi
fi

if [ -n "$2" ] && dpkg --compare-versions "$2" lt 1.36-42; then
	update-rc.d -f univention-maintenance remove || true
fi

eval $(univention-config-registry shell server/role)

if [ "$server_role" = "domaincontroller_master" ] || [ "$server_role" = "domaincontroller_backup" ]; then
	/usr/lib/univention-install/41univention-updater.inst || true
fi

if [ ! -d /etc/apt/sources.list.d ]; then
	mkdir -p /etc/apt/sources.list.d
fi

if [ -n "$2" ] && dpkg --compare-versions "$2" lt 3.1.1-5; then
	#create an empty sources.list
	if [ -e "/etc/apt/sources.list" ]; then
		mv /etc/apt/sources.list /etc/apt/sources.list.unused
		touch /etc/apt/sources.list
	fi
fi

# the UCR debhelper section will be added before the pycentral section,
# so we added the pycentral section manually
if which pycentral >/dev/null 2>&1; then
	pycentral pkginstall univention-updater
fi

#DEBHELPER#

if [ -x /etc/init.d/univention-management-console-server ]; then
	/etc/init.d/univention-management-console-server crestart
fi

# activate the online repository, initialize mirror configuration
univention-config-registry set	repository/online?yes \
								repository/online/port?80 \
								repository/online/maintained?yes \
								repository/online/unmaintained?no \
								repository/online/hotfixes?no \
								repository/mirror/threads?10 \
								repository/mirror/basepath?/var/lib/univention-repository

# set the repository server
eval $(univention-config-registry shell hostname domainname update/server)
if [ -n "$update_server" -a "$update_server" != "http://download.univention.de" ]; then
	update_server=`echo "$update_server" | sed -e 's|http://||' | sed -e 's|/.*||'`
	univention-config-registry set repository/online/server?"$update_server"
elif host univention-repository.$domainname &> /dev/null; then
	univention-config-registry set repository/online/server?"univention-repository.$domainname"
else
	univention-config-registry set repository/online/server?"apt.univention.de"
fi

# unset old variables
univention-config-registry unset update/server update/security/server update/security/prefix update/directory

# define update commands
univention-config-registry set \
		update/commands/distupgrade?"apt-get -o DPkg::Options::=--force-confold -y --force-yes -u dist-upgrade" \
		update/commands/update?"apt-get update" \
		update/commands/show?"apt-get show" \
		update/commands/install?"apt-get -o DPkg::Options::=--force-confold -y --force-yes install" \
		update/commands/remove?"apt-get --yes remove" \
		update/commands/configure?"dpkg --configure -a" \
		update/commands/distupgrade/simulate?"apt-get -o DPkg::Options::=--force-confold -y --force-yes -us dist-upgrade" \
		update/commands/upgrade/simulate?"apt-get -o DPkg::Options::=--force-confold -y --force-yes -us upgrade"

exit 0
