# base/univention-updater/debian/univention-updater.cron.d
PATH=/usr/sbin:/usr/bin:/sbin:/bin

@reboot root /usr/share/univention-updater/set-apache-permissions
