# base/univention-updater/debian/univention-updater.cron.d
@reboot root /usr/share/univention-updater/set-apache-permissions
