# management/univention-directory-listener/conffiles/etc/systemd/system/univention-directory-listener.service.d/20_start_order.conf
@%@UCRWARNING=# @%@
[Unit]
@!@
server_role = configRegistry.get('server/role', '')
if server_role == 'domaincontroller_master':
	print('After=univention-directory-notifier.service')
else:
	print('Before=univention-directory-notifier.service')
@!@
