# services/univention-shares/conffiles/cron.d/univention-shares
# Run queue every 15 minutes
05,20,35,50 *     * * *     /usr/sbin/univention-share-replication
