# services/univention-scalix/debian/univention-scalix.conffiles
/etc/univention/apache/95scalix-jakarta.conf
/etc/univention/templates/files/etc/init.d/scalix-tomcat
