# services/univention-scalix/debian/univention-scalix.conffiles
/etc/univention/templates/files/etc/init.d/scalix-tomcat
/etc/univention/templates/files/etc/opt/scalix-tomcat/connector/jk/workers.conf
