# services/univention-scalix/debian/univention-scalix-tomcat.conffiles
/etc/univention/templates/files/etc/init.d/scalix-tomcat
