# services/univention-scalix/debian/univention-scalix.conffiles
/etc/apache2/conf.d/scalix-jakarta.conf
/etc/univention/templates/files/etc/init.d/scalix-tomcat
