# saml/univention-saml/systemd/system/stunnel4.service.d/10-watchdog.conf
[Service]
RemainAfterExit=no
Restart=always

