# desktop/univention-gdm/conffiles/etc/init.d/gdm
#! /bin/sh
@%@UCRWARNING=# @%@

### BEGIN INIT INFO
# Provides:          x-display-manager gdm
# Should-Start:      console-screen acpid
# Required-Start:    $local_fs $remote_fs
# Required-Stop:     $local_fs $remote_fs
# Default-Start:     2 3 4 5
# Default-Stop:      S 0 1 6
# Short-Description: GNOME Display Manager
# Description:       Debian init script for the GNOME Display Manager
### END INIT INFO
#
# Author:	Ryan Murray <rmurray@debian.org>
#

# Copyright (C) 2004-2012 Univention GmbH
#
# http://www.univention.de/
#
# All rights reserved.
#
# The source code of the software contained in this package
# as well as the source package itself are made available
# under the terms of the GNU Affero General Public License version 3
# (GNU AGPL V3) as published by the Free Software Foundation.
# 
# Binary versions of this package provided by Univention to you as
# well as other copyrighted, protected or trademarked materials like
# Logos, graphics, fonts, specific documentations and configurations,
# cryptographic keys etc. are subject to a license agreement between
# you and Univention and not subject to the GNU AGPL V3.
# 
# In the case you use the software under the terms of the GNU AGPL V3,
# the program is provided in the hope that it will be useful,
# but WITHOUT ANY WARRANTY; without even the implied warranty of
# MERCHANTABILITY or FITNESS FOR A PARTICULAR PURPOSE. See the
# GNU Affero General Public License for more details.
# 
# You should have received a copy of the GNU Affero General Public
# License with the Debian GNU/Linux or Univention distribution in file
# /usr/share/common-licenses/AGPL-3; if not, see
# <http://www.gnu.org/licenses/>.

# do not start the GDM if the setup tools will be started
if [ -x /usr/sbin/univention-config-registry ]; then
	eval "$(/usr/sbin/univention-config-registry shell system/setup/boot/start)"
	if [ -n "$system_setup_boot_start" ] && [ "$system_setup_boot_start" = "yes" ] || [ "$system_setup_boot_start" = "true" ]; then
		exit 0
	fi

	max ()
	{
		if [ "$1" -gt "$2" ]; then
			echo -n "$1"
		else
			echo -n "$2"
		fi
	}

	# calculate xorg/display/virtual if xorg/resolution/secondary is set
	eval "$(/usr/sbin/univention-config-registry shell xorg/resolution/secondary xorg/display/relative-position xorg/display/virtual xorg/resolution)"
	if [ -n "$xorg_resolution_secondary" ]; then
		x_prim="${xorg_resolution%x*}"
		y_prim="${xorg_resolution##*x}"
		x_sec="${xorg_resolution_secondary%x*}"
		y_sec="${xorg_resolution_secondary##*x}"

		if [ -n "$xorg_display_relative_position" ]; then
			if [ "$xorg_display_relative_position" = "right" -o "$xorg_display_relative_position" = "left" ]; then
				new_value="$((x_prim+x_sec)) $(max $y_prim $y_sec)"
			elif [ "$xorg_display_relative_position" = "above" -o "$xorg_display_relative_position" = "below" ]; then
				new_value="$(max $x_prim $x_sec) $((y_prim+y_sec))"
			fi
		fi
		if [ -n "$new_value" ] && [ "$new_value" != "$xorg_display_virtual" ]; then
			if [ -z "$xorg_ignore_display_virtual" ]; then
				univention-config-registry set xorg/display/virtual="$new_value"
			fi
		fi
	fi
fi

set -e

PATH=/sbin:/bin:/usr/sbin:/usr/bin
DAEMON=/usr/sbin/gdm
DEBCONF_DAEMON=/usr/sbin/gdm

test -x $DAEMON || exit 0

if [ -r /etc/default/locale ]; then
  . /etc/default/locale
  export LANG LANGUAGE
fi

. /lib/lsb/init-functions

# To start gdm even if it is not the default display manager, change
# HEED_DEFAULT_DISPLAY_MANAGER to "false."
HEED_DEFAULT_DISPLAY_MANAGER=true
DEFAULT_DISPLAY_MANAGER_FILE=/etc/X11/default-display-manager

case "$1" in
  start)
	# check ucr autostart setting
	if [ -f "/usr/share/univention-config-registry/init-autostart.lib" ]; then
	    . "/usr/share/univention-config-registry/init-autostart.lib"
	    check_autostart gdm gdm/autostart
	fi

	# check for autostart scripts
	if [ -x /usr/sbin/univention-config-registry ]; then
		eval "$(/usr/sbin/univention-config-registry shell univentionAutoStartScript gdm/autostart/user)"
		if [ -n "$univentionAutoStartScript" ] && [ "$univentionAutoStartScript" != "none" -a "$univentionAutoStartScript" != "None" ]; then
			/usr/share/univention-gdm-sessions/autostart >>/var/log/x-autostart.log 2>&1 &
			exit 0
		fi
	fi

  	if [ -e "$DEFAULT_DISPLAY_MANAGER_FILE" -a "$HEED_DEFAULT_DISPLAY_MANAGER" = "true" -a "$(cat $DEFAULT_DISPLAY_MANAGER_FILE 2> /dev/null)" != "$DEBCONF_DAEMON" ]; then
		log_action_msg "Not starting GNOME Display Manager; it is not the default display manager"
	else
		# don't show log messages during boot, otherwise we see a console with the gdm starting message
		! pidof usplash > /dev/null && log_daemon_msg "Starting GNOME Display Manager" "gdm"
		start_daemon $DAEMON; result=$?
		! pidof usplash > /dev/null && log_end_msg $result
	fi
  ;;
  stop)
	log_daemon_msg "Stopping GNOME Display Manager" "gdm"
	set +e
	start-stop-daemon --stop --quiet --pidfile /var/run/gdm.pid \
		--name gdm --retry 5
	set -e
	log_end_msg $?
  ;;
  reload)
  	log_daemon_msg "Scheduling reload of GNOME Display Manager configuration" "gdm"
	set +e
	start-stop-daemon --stop --signal USR1 --quiet --pidfile \
		/var/run/gdm.pid --name gdm
	set -e
	log_end_msg $?
  ;;
  restart|force-reload)
	$0 stop
	$0 start
  ;;
  *)
	echo "Usage: /etc/init.d/gdm {start|stop|restart|reload|force-reload}"
	exit 1
  ;;
esac

exit 0
