# services/univention-bind/conffiles/etc/runit/univention-bind-proxy/run
#!/bin/sh

PATH=/sbin:/bin:/usr/sbin:/usr/bin

OPTS="-c /etc/bind/named.conf.proxy -u bind -f"

exec /usr/sbin/named $OPTS
