# services/univention-bind/conffiles/etc/runit/univention-bind-proxy/run
#!/bin/sh

PATH=/sbin:/bin:/usr/sbin:/usr/bin

OPTS="-c /etc/bind/named.conf.proxy -u bind -f"

# check ucr autostart setting
if [ -f "/usr/share/univention-config-registry/init-autostart.lib" ]; then
    source "/usr/share/univention-config-registry/init-autostart.lib"
    check_autostart univention-bind-proxy bind/autostart
fi

exec /usr/sbin/named $OPTS
