# services/univention-bind/conffiles/etc/runit/univention-bind-proxy/run
#!/bin/sh

PATH=/sbin:/bin:/usr/sbin:/usr/bin

OPTS="-c /etc/bind/named.conf.proxy -u bind -f"

eval `univention-config-registry shell bind/autostart`
if [ -n "$bind_autostart" ]; then
	if [ "$bind_autostart"  == "no" ] || [ "$bind_autostart"  == "false" ]; then
		echo " bind is disabled"
		exit 0
	fi
fi


exec /usr/sbin/named $OPTS
