# saml/univention-saml/systemd/system/univention-stunnel4-watchdog.service
[Unit]
Description=Restarting stunnel4 if it hangs (Bug 52196)

[Service]
ExecStart=/usr/lib/univention-saml/univention-stunnel4-watchdog start

[Install]
WantedBy=multi-user.target


