# services/univention-dhcp/conffiles/etc/runit/univention-dhcp/run
#! /bin/sh -e

eval `/usr/sbin/univention-config-registry shell dhcpd/autostart`

if [ "$dhcpd_autostart" = "yes" ] ; then
	exec /usr/sbin/dhcpd3 -q -f >>/var/log/daemon.log
else
	logger -p local0.notice -t UNIVENTION-DHCP "dhcpd disabled in configuration registry - not started."
fi
