# services/univention-dhcp/conffiles/etc/runit/univention-dhcp/run
#! /bin/sh -e


exec /usr/sbin/dhcpd3 -q -f >>/var/log/daemon.log
