# management/univention-management-console/conffiles/etc/runit/univention-management-console-server/run
#!/bin/sh

# Warning: This file is auto-generated and might be overwritten by
#          univention-baseconfig.
#          Please edit the following file instead:
# Warnung: Diese Datei wurde automatisch generiert und kann durch
#          univention-baseconfig überschrieben werden.
#          Bitte bearbeiten Sie an Stelle dessen die folgende Datei:
#
# 	/etc/univention/templates/files/etc/runit/univention-management-console-server/run
#


PATH=/sbin:/bin:/usr/sbin:/usr/bin

@!@
cmd = 'exec /usr/sbin/univention-management-console-server -n -d @%@umc/server/debug/level@%@ %s >>/var/log/univention/management-console-server.log 2>&1'
lang = ''
if baseConfig.has_key('umc/web/language'):
	lang = '-l %s' % baseConfig['umc/web/language']

print cmd % lang
@!@
